# knowledger-processor.service
[Unit]
Description=Knowledger Processor

[Service]
Type=simple
WorkingDirectory=/root/knowledger-miner/src/processor
ExecStart=/root/knowledger-miner/knowledger-processor.sh
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
